# certbot.renew.service.template
[Unit]
Description=Service for snap application certbot.renew

[Service]
ExecStart=${NGINX_ROOT_FOLDER}/renew-cert.sh