# certbot.renew.timer
[Unit]
Description=Timer for renewing certbot certifactes

[Timer]
Unit=certbot.renew.service
OnCalendar=daily

[Install]
WantedBy=timers.target